# database/backup/database_backup.timer
[Unit]
Description=Daily timer to run database_backup.service at 19:00

[Timer]
# Run every day at 19:00 (7:00 PM)
OnCalendar=*-*-* 19:00:00
# If the machine was off at the scheduled time, run once at next boot
Persistent=true
# The service unit to activate
Unit=database_backup.service

[Install]
WantedBy=timers.target